# baselines/vision_obstacle_avoidance_legacy/src/vision_obstacle_avoidance/systemd/voice-daemon.service
[Unit]
Description=Vision Obstacle Avoidance — Voice Dictation Daemon
Documentation=https://github.com/user/vision_obstacle_avoidance
After=graphical-session.target sound.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/bash -c 'test -d /dev/shm || mkdir -p /tmp/.voice-daemon'
ExecStart=__PROJECT_DIR__/.venv/bin/python3 __PROJECT_DIR__/voice_daemon.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Desktop environment passthrough
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=HF_ENDPOINT=https://hf-mirror.com

# Security hardening
PrivateTmp=true
NoNewPrivileges=true

# 语音守护进程不适用 watchdog（大部分时间阻塞在信号等待上）
RestartSec=10

[Install]
WantedBy=default.target
